# spotifyd-display.service
[Unit]
Description=Spotifyd E-Paper Display Service
After=network.target sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/spotifyd-epaper-display
Environment="PYTHONUNBUFFERED=1"
Environment="DISPLAY=:0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
ExecStart=/usr/bin/python3 /home/pi/spotifyd-epaper-display/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
